# deploy/roles/app/templates/phantom.conf.j2
start on started networking
stop on shutdown
exec /usr/bin/phantomjs --webdriver=127.0.0.1:8910
